# rootfs/etc/systemd/system/altitude-systemd-audit.timer
[Unit]
Description=Refresh Altitude systemd audit table

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=altitude-systemd-audit.service

[Install]
WantedBy=timers.target
